# scripts/upstart/hulk.conf
description "hulk - Haskell IRC server"

start on filesystem
stop on runlevel [016]

respawn

# Change according to installed location
exec hulk -c /path/to/hulk.conf
